# supervisor/guacamole.sv.conf
[program:guacamole]
command=/bin/bash -c "service guacd start && service tomcat7 restart"
redirect_stderr=true